# deploy/cemini/systemd/cemini-wc-cross-venue.service
[Unit]
Description=World Cup Bot Module 6 — PM↔Kalshi cross-venue alert loop
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/cemini
EnvironmentFile=-/opt/cemini/.env
Environment=WC_DRY_RUN=true
ExecStart=/bin/bash -c '/opt/cemini/scripts/wc_run.sh cross-venue-scan --loop --alert-only 2>&1 | while IFS= read -r line; do echo "$(date -u +%Y-%m-%dT%H:%M:%SZ) $line"; echo "$(date -u +%Y-%m-%dT%H:%M:%SZ) $line" >> /opt/cemini/logs/wc_cross_venue_alerts.jsonl; [[ "$line" == ALERT* ]] && logger -t cemini-wc "$line"; done'
Restart=always
RestartSec=30
User=root

[Install]
WantedBy=multi-user.target
