# deploy/cemini/systemd/cemini-wc-calendar-guard.service
[Unit]
Description=World Cup Bot calendar cancel-window check
After=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/cemini
EnvironmentFile=-/opt/cemini/.env
Environment=WC_DRY_RUN=true
ExecStart=/bin/bash -c '/opt/cemini/scripts/wc_run.sh calendar --cancel-window --min-hours 10 >> /opt/cemini/logs/cron_wc_calendar.log 2>&1'
User=root
