# deploy/cemini/systemd/cemini-wc-preflight.service
[Unit]
Description=World Cup Bot preflight (geoblock + CLOB auth probe)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/cemini
EnvironmentFile=-/opt/cemini/.env
EnvironmentFile=-/opt/cemini/.env-polymarket
Environment=WC_DRY_RUN=true
ExecStart=/bin/bash -c '/opt/cemini/scripts/wc_run.sh preflight >> /opt/cemini/logs/cron_wc_preflight.log 2>&1'
User=root
