# deploy/cemini/systemd/cemini-wc-live-plan.service
[Unit]
Description=World Cup Bot live plan --record (operator GO required)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/cemini
EnvironmentFile=-/opt/cemini/.env
EnvironmentFile=-/opt/cemini/.env-polymarket
Environment=WC_DRY_RUN=false
Environment=WC_LEDGER_PATH=/opt/cemini/logs/wc_ledger.jsonl
ExecStart=/bin/bash -c '/opt/cemini/scripts/wc_run.sh preflight && /opt/cemini/scripts/wc_run.sh plan --record >> /opt/cemini/logs/cron_wc_live_plan.log 2>&1'
User=root
