# deploy/systemd/units/world-cup-bot-shadow-plan.service
[Unit]
Description=World Cup Bot shadow plan --record (DRY_RUN)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=@INSTALL_ROOT@
EnvironmentFile=-@INSTALL_ROOT@/.env
Environment=WC_DRY_RUN=true
Environment=WC_LEDGER_PATH=@INSTALL_ROOT@/data/local/shadow_ledger.jsonl
ExecStart=/bin/bash -c '@INSTALL_ROOT@/bin/wc_run.sh plan --record --liquidity-gate >> @INSTALL_ROOT@/logs/cron_shadow_plan.log 2>&1'
User=root
