# deploy/systemd/units/world-cup-bot-cross-venue-reconcile.timer
[Unit]
Description=World Cup Bot weekly cross-venue reconcile (Sun 08:00 UTC)

[Timer]
OnCalendar=Sun *-*-* 08:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
