# deploy/systemd/units/world-cup-bot-tournament-ops.timer
[Unit]
Description=World Cup Bot daily tournament ops check

[Timer]
OnCalendar=*-*-* 06:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
