# deploy/systemd/units/world-cup-bot-rewards-sync.service
[Unit]
Description=World Cup Bot daily rewards sync (L2 creds required — Phase 2+)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=@INSTALL_ROOT@
EnvironmentFile=-@INSTALL_ROOT@/.env
EnvironmentFile=-@INSTALL_ROOT@/.env.trading
Environment=WC_DRY_RUN=true
Environment=WC_LEDGER_PATH=@INSTALL_ROOT@/data/local/shadow_ledger.jsonl
ExecStart=/bin/bash -c '@INSTALL_ROOT@/bin/wc_run.sh rewards sync --record >> @INSTALL_ROOT@/logs/cron_rewards.log 2>&1'
User=root
